# .config/systemd/user/app@.service
[Unit]
Description=Application '%i'

[Service]
ExecStart=/usr/bin/env %i
