# .config/systemd/user/background@.service
[Unit]
Description=Set Background In Apps To '%i'
PartOf=background.target
Before=background.target
After=graphical.target

[Service]
ExecStart=/usr/bin/systemctl --user set-environment BACKGROUND=%i
ExecStartPost=-/bin/sh -c '[ "${BACKGROUND}" != "%i" ] && \
    exec systemctl --user stop background@${BACKGROUND}'
Type=oneshot
RemainAfterExit=yes

[Install]
RequiredBy=background.target
DefaultInstance=light
