# .config/systemd/user/background.target
[Unit]
Description=X Background Configured

[Install]
WantedBy=workplace.target
RequiredBy=display-manager.target
Also=background@.service
