# .config/systemd/user/background.target
[Unit]
Description=X Background Configured
BindsTo=display-manager.target
