# .config/systemd/user/display-manager.target
[Unit]
Description=Window manager online (i3)
Requires=graphical.target
After=graphical.target

[Install]
WantedBy=default.target
