# .config/systemd/user/environment.target
[Unit]
Description=User Environment Online
PartOf=workplace.target

[Install]
RequiredBy=workplace.target
