# .config/systemd/user/environment.target
[Unit]
Description=User Environment Online

[Install]
RequiredBy=workplace.target
