# .config/systemd/user/environment.target
[Unit]
Description=User Environment Online

[Install]
Also=environment.service
