# .config/systemd/user/keyboard.target
[Unit]
Description=Fully Configured Keyboard

[Install]
WantedBy=workplace.target
Also=keyboard.service
Also=sxhkd.service
Also=keynav.service
