# .config/systemd/user/multiplexer.target
[Unit]
Description=Terminal Multiplexer

[Install]
Also=tmux.service
