# .config/systemd/user/multiplexer.target
[Unit]
Description=Terminal Multiplexer

[Install]
WantedBy=workplace.target
Also=tmux.service
