# .config/systemd/user/tmux.service
[Unit]
Description=Start Tmux Server
Before=multiplexer.target

[Service]
Type=forking
;ExecStart=/usr/bin/tmux new -ds __systemd sh
ExecStart=/usr/bin/tmux start-server
ExecStop=/usr/bin/tmux kill-server

[Install]
WantedBy=multiplexer.target
