# .config/systemd/user/vimruntime-init.service
[Unit]
Description="vimruntime-init"

[Service]
ExecStart=/usr/bin/sh -c vimruntime-init
Type=oneshot

[Install]
WantedBy=environment.target
