# .config/systemd/user/vm@.service
[Unit]
Description=Start VirtualBox Virtual Machine
Requires=graphical.target

[Service]
ExecStart=/usr/bin/vboxheadless --startvm %i --vrde on
