# .config/systemd/user/xclipboard-synchronizer.service
[Unit]
Description=Synchronize various X clipboard buffers for sane result
Requires=graphical.target
After=graphical.target
Before=display-manager.target

[Service]
ExecStart=/bin/sh clipboard-sync

[Install]
WantedBy=display-manager.target
