# .config/systemd/user/xclipboard-synchronizer.service
[Unit]
Description=Synchronize various X clipboard buffers for sane result
Requires=graphical.target
After=graphical.target

[Service]
ExecStart=/bin/xcsyncd --mode CLIP_TO_PRIM

[Install]
WantedBy=workplace.target
