# .config/systemd/user/xclipboard-synchronizer.service
[Unit]
Description=Synchronize various X clipboard buffers for sane result
Requires=graphical.target
After=graphical.target
Before=display-manager.target

[Service]
Type=forking
ExecStart=/bin/bash -c ' \
    autocutsel -fork -pause 150 -selection PRIMARY && \
    autocutsel -fork -pause 150 -selection CLIPBOARD \
'

[Install]
WantedBy=display-manager.target
