# <root>/etc/udev/rules.d/99-external-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{ID_PATH}=="pci-0000:00:02.0", RUN+="/usr/bin/su operator -lc 'DISPLAY=:0 external-monitor'"
