# <root>/etc/udev/rules.d/99-reconfig-kbd.rules
ACTION=="add", SUBSYSTEMS=="usb", SUBSYSTEM=="input", ENV{ID_INPUT_KEYBOARD}=="1", ENV{REL}=="[0-9]*", RUN+="/home/operator/bin/reconfig-kbd"
