# _/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
TTYVTDisallocate=no
ExecStart=
ExecStart=-/usr/bin/agetty --autologin operator --noclear %I $TERM
