# _/usr/lib/systemd/system/zeusd.timer
[Unit]
Description=zeusd timer: run nightly backups

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
