# rootfs/etc/systemd/system/user@.service.d/00-add-unit-path.conf
[Service]
Environment=SYSTEMD_UNIT_PATH=/home/s.seletskiy/.systemd/units/:
